# lamp/supervisord-postfix.conf
[program:postfix]
command=/etc/init.d/postfix start
numprocs=1
autostart=true
autorestart=true