# lamp/supervisord-postfix.conf
[program:postfix]
process_name=master
directory=/etc/postfix
command=/usr/sbin/postfix -c /etc/postfix start
autostart=true
autorestart=true

